# scripts/crontab
*/11 * * * * cd $HOME/mirage-decks-deployment && git pull
